# roles/tshark-listener/templates/rc.local.j2
#!/bin/bash

ip link set {{ tshark_listener_int }} promisc on
ip link set {{ tshark_listener_int }} up


exit 0